# archives/filesystem/common/etc/update-motd.d/10-help-text
#!/bin/sh

printf "\n"
printf " * Website:        https://www.namtso.com\n"
printf "\n"
